# extras/systemd/docker_pipe.service
[Service]
Type=simple
ExecStart=/usr/local/bin/docker_pipe.sh
PIDFile=/tmp/docker_pipe.pid
